# test/crontab/carrot
MAILTO=carrot@carrot.org

29	10	1   1	0 	/home/carrot/logrotate.sh --from 2017/02/20
